# wave-report-server/com.claude-ops.wave-report.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-ops.wave-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>$HOME/.bun/bin/bun</string>
        <string>run</string>
        <string>$HOME/.claude-ops/wave-report-server/server.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>$HOME/.claude-ops/wave-report-server</string>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>$HOME/Library/Logs/wave-report-server.log</string>
    <key>StandardErrorPath</key>
    <string>$HOME/Library/Logs/wave-report-server.log</string>
</dict>
</plist>
